# deploy/systemd/aquapi-collect.service
[Unit]
Description=aquapi temperature collector
Documentation=https://github.com/youhey/aquapi
After=local-fs.target

[Service]
Type=simple
User=aquapi
Group=aquapi
WorkingDirectory=/opt/aquapi
Environment=AQUAPI_CONFIG=/etc/aquapi/aquapi.json
ExecStart=/opt/aquapi/.venv/bin/python -m aquapi.cli collect --config ${AQUAPI_CONFIG}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

